# dot-config/systemd/user/xdg-desktop-portal.service.d/override.conf
[Unit]
; Requisite=
; Requires=graphical-session.target
